# deb/tink-backend-webhook.service
[Unit]
Description=tink-backend-webhook

[Service]
User=tink-backend-webhook
Group=tink-backend-webhook
KillSignal=SIGINT
TimeoutStopSec=20
Restart=on-failure
RestartSec=10
Environment=MAX_HEAP_SIZE=256m
EnvironmentFile=-/etc/default/tink-backend-webhook
WorkingDirectory=/usr/share/tink-backend-webhook

ExecStart=/usr/bin/java \
          -Xmx${MAX_HEAP_SIZE} \
          -Xloggc:/var/log/tink/webhook/gc.log \
          -XX:+UseGCLogFileRotation \
          -XX:NumberOfGCLogFiles=5 \
          -XX:GCLogFileSize=512K \
          -Dnetworkaddress.cache.ttl=60 \
          -XX:-OmitStackTraceInFastThrow \
          -XX:+HeapDumpOnOutOfMemoryError \
          -XX:HeapDumpPath=/var/tmp/tink/heapdumps/tink-backend-webhook \
          -XX:+ExitOnOutOfMemoryError \
          -jar /usr/share/tink-backend-webhook/webhook-service.jar \
          server /etc/tink/webhook-server.yml

[Install]
WantedBy=multi-user.target
